# autostart/tiosk-mono-audio.desktop
[Desktop Entry]
Type=Application
Name=TIOSK Mono Audio
Exec=/home/kiosk/tiosk_mono_audio.sh
Terminal=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
