# autostart/touch-calibrate.desktop
[Desktop Entry]
Type=Application
Name=Touch Calibration
Exec=/usr/local/bin/calibrate-touch.sh
Terminal=false
X-GNOME-Autostart-enabled=true
