# autostart/xfce4-screensaver.desktop
[Desktop Entry]
Type=Application
Name=xfce4-screensaver
Exec=/bin/true
Hidden=true
X-GNOME-Autostart-enabled=false
